# waterme/rootfs/etc/services.d/waterme/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting WaterME Irrigation Controller..."

# Set workdir to /app where our code lives
cd /app

# Ensure we have a place for data
mkdir -p /data

bashio::log.info "Launching FastAPI Backend..."

# Run uvicorn
# --host 0.0.0.0 to bind to all interfaces (required for Ingress/Docker networking)
# --port 8099 matches our config.yaml ingress_port
exec uvicorn backend.server:app --host 0.0.0.0 --port 8099 --log-level info
